# systemd/ape-signal-tick@.service
[Unit]
Description=Ape Signal paper-trading tick (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ape-signal
EnvironmentFile=/etc/ape-signal.env
# HOME so the Claude CLI finds the interactive `claude login` credentials in /root/.claude.
Environment=HOME=/root
ExecStart=/usr/bin/npx tsx src/paper/tick.ts %i
