# systemd/ape-signal-tick.timer
[Unit]
Description=Ape Signal paper-trading monitor tick (every 5 min 15:30-21:55 Europe/Berlin)

[Timer]
# US session in Berlin wall-clock, 5-minute monitor raster (ADR 0003). The
# 22:00 closing tick is a separate timer (ape-signal-tick-close.timer).
# No Persistent: a missed tick is stale; the next one is at most 5 minutes away.
OnCalendar=Mon..Fri *-*-* 15:30/5:00 Europe/Berlin
OnCalendar=Mon..Fri *-*-* 16..21:00/5:00 Europe/Berlin
Unit=ape-signal-tick@Tick.service

[Install]
WantedBy=timers.target
